# templates/default/zmq_broker.conf.erb
# zmq_broker - ZeroMQ Broker Application
#
# zmq_broker receives messages on port 5555 and publishes them on port 5556

description "ZeroMQ Broker"

start on net-device-up
stop on runlevel [!2345]

respawn
respawn limit 10 5

setuid zmq_broker
setgid zmq_broker

pre-start script
    getent passwd <%= @user %> >/dev/null 2>&1 || { stop; exit 0; }
    test -x <%= @env['HOME'] %>/zmq_broker_go || { stop; exit 0; }
end script

script
  cd <%= @env['HOME'] %>
  exec ./zmq_broker_go <%= @flags %> >> /var/log/zmq_broker.log 2>&1
end script
